# unit-files/backup@.timer
[Unit]
Description=Hourly backup of www and mysql content to a git repo

[Timer]
OnUnitActiveSec=1h

[Install]
WantedBy=timers.target
