# unit-files/backup@.timer
[Unit]
Description=Hourly backup of www and mysql content to a git repo

[Timer]
OnActiveSec=20
OnUnitActiveSec=60min

