# unit-files/git-puller.timer
[Unit]
Description=Git pull every 5 minutes

[Timer]
OnUnitActiveSec=5min

