# unit-files/git-puller.timer
[Unit]
Description=Git pull every 5 minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=5min

[Install]
WantedBy=multi-user.target

