# unit-files/gitpuller@.service
[Unit]
Description= git puller

[Service]
Environment=USER=%i
ExecStart=/data/indiehosters/importers/gitpuller.sh
