# unit-files/haproxy.path
[Path]
PathExists=/data/server-wide/haproxy/haproxy.cfg

[Install]
WantedBy=multi-user.target
