# unit-files/haproxy.path
[Path]
PathExists=/data/runtime/haproxy/haproxy.cfg

[Install]
WantedBy=multi-user.target

