# unit-files/nginx-gitpuller@.path
[Unit]
Description=Start the timer if there is a GITURL file

# Dependency binding
BindsTo=nginx@%i.service

[Path]
PathExists=/data/domains/%i/nginx/data/GITURL
Unit=nginx-gitpuller@%i.timer
