# unit-files/nginx-gitpuller@.service
[Unit]
Description= git puller

[Service]
Environment=USER=%i
Environment=APP=nginx
ExecStart=/data/indiehosters/importers/gitpuller.sh
