# unit-files/nginx-gitpuller@.timer
[Unit]
Description=Run git puller every 10 minutes

[Timer]
OnUnitActiveSec=10min

[Install]
WantedBy=timers.target
