# unit-files/nginx-gitpuller@.path
[Path]
PathExists=/data/per-user/%i/nginx/data/git-url.txt
