# unit-files/ocsp.timer
[Unit]
Description=Daily timer for OCSP stapling

[Timer]
OnUnitActiveSec=1day

