# unit-files/ocsp.timer
[Unit]
Description=Daily timer for OCSP stapling

[Timer]
OnBootSec=15min
OnUnitActiveSec=1day

[Install]
WantedBy=multi-user.target

