# unit-files/static-importer@.service
[Unit]
Description=static importer

# Dependency ordering
Before=static@%i.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/bin/bash -euxc ' \
  domain_folder=/data/domains/%i; \
  cp $domain_folder/TLS/%i.pem /data/runtime/haproxy/approved-certs/%i.pem; \
  folder=$domain_folder/static/www-content; \
  mkdir -p $folder; \
  if [ ! -f $folder/index.html ]; then \
    echo Hello %i > $folder/index.html; \
  fi;'

[Install]
WantedBy=static@%i.service
