# unit-files/swapon.service
[Unit]
Description=Turn on swap

[Service]
Type=oneshot
ExecStart=/sbin/swapon /swap

[Install]
WantedBy=local.target
