# unit-files/wordpress-subdir-gitpuller@.service
[Unit]
Description= git puller

[Service]
Environment=USER=%i
Environment=APP=wordpress-subdir
ExecStart=/data/indiehosters/importers/gitpuller.sh
