# unit-files/wordpress-subdir-gitpuller@.path
[Path]
PathExists=/data/per-user/%i/wordpress-subdir/data/GITURL
